# scripts/launchd/dev.skz.ccmux.traffic.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Daily run of scripts/traffic.sh — snapshots GitHub traffic for
  skzv/ccmux + skzv/homebrew-tap and the latest release's download
  count into ~/.local/state/ccmux/traffic.csv. The script needs `gh`
  authenticated with push access to both repos.

  Install:
    cp scripts/launchd/dev.skz.ccmux.traffic.plist ~/Library/LaunchAgents/
    launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/dev.skz.ccmux.traffic.plist
    launchctl kickstart gui/$(id -u)/dev.skz.ccmux.traffic   # fire once now

  Uninstall:
    launchctl bootout gui/$(id -u)/dev.skz.ccmux.traffic
    rm ~/Library/LaunchAgents/dev.skz.ccmux.traffic.plist

  Assumes the ccmux clone lives at $HOME/Projects/ccmux. Edit the
  ProgramArguments string below if yours is somewhere else.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>dev.skz.ccmux.traffic</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-lc</string>
    <string>bash "$HOME/Projects/ccmux/scripts/traffic.sh"</string>
  </array>

  <!-- Daily at 09:00 local time. -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key><integer>9</integer>
    <key>Minute</key><integer>0</integer>
  </dict>

  <!-- Run-on-load so a fresh install captures a row immediately,
       not 24 h later. -->
  <key>RunAtLoad</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/tmp/ccmux-traffic.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/ccmux-traffic.err.log</string>
</dict>
</plist>
